# dot_local/bin/docker-pre-shutdown.service
[Unit]
Description=Scale down docker-swarm services before system shutdown
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/true
ExecStop=/usr/local/bin/docker-pre-shutdown.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

